# deployment/supervisor/queue.conf
[program:laraq]
process_name=%(program_name)s_%(process_num)02d
command=php /var/www/backend.your-domain.com/artisan queue:work --queue=send_otp --sleep=3 --tries=3
autostart=true
autorestart=true
user=ubuntu
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/backend.your-domain.com/storage/logs/queue-worker.log
